# etc/systemd/system/rfid.service
[Unit]
Description=RFID Reader

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/rfid.py

[Install]
WantedBy=multi-user.target
